# com.local.shanghai-weather-email.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local.shanghai-weather-email</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/ABSOLUTE/PATH/TO/weather_email/run_daily_weather.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>16</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/ABSOLUTE/PATH/TO/weather_email/weather_email.log</string>
    <key>StandardErrorPath</key>
    <string>/ABSOLUTE/PATH/TO/weather_email/weather_email.err.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
